# scapy-vxlan/scapy-vxlan.service
[Unit]
Description=Spark service

[Service]
ExecStart=/usr/local/bin/scapy-vxlan.py

[Install]
WantedBy=multi-user.target
